# tier0/desktop/systemd/packagekit.service.d/override.conf
[Unit]
# Remove ConditionPathExists directive added by packagekit.service, because we now have our own agent!
ConditionPathExists=